# etc/cron.d/Sauv2Stephdl.dyndns.org
*/30 * * * * root rsync -avz --progress --delete /home/e-smith/files/users/ root@stephdl.dyndns.org:/home/e-smith/files/users/
