# etc/cron.d/sauvMysql2stephdl.dyndns.org
7 7 * * * root rsync -avz --delete /root/backup/db/ root@stephdl.dyndns.org:/root/SauvMysqlAubrac
